# etc/crontab
* * * * * /bin/echo "test" > /var/log/test.log
